# Systemd/change_wallpaper_random.timer
[Unit]
Description=Change Wallpaper every 10 minutes
After=graphical-session.target

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
